# live_trend.cron
30 9 * * 1-6cd /app && /app/build/bin/Release/live_trend >> /var/log/cron.log 2>&1
