# live_trend.cron
30 9 * * 1-6 sh /app/scripts/run_live_trend.sh >> /var/log/cron.log 2>&1
